# vpn-passthrough/dummy-app/dummy-app.service
[Unit]
Description=Dummy App

[Service]
Type=simple
ExecStart=dummy-app serve