# vpn-passthrough/dummy-app/dummy-app@.service
[Unit]
Description=Dummy App for %I

[Service]
Type=simple
ExecStart=dummy-app --instance "%I" serve