# vpn-passthrough/vpn-passthrough/vpn-passthrough@.service
[Unit]
Description=VPN Passthrough for %Install
PartOf=vpn-passthrough@%i.service

[Service]
ExecStart=oneshot
RemainsAfterExit=yes