# vpn-passthrough/vpn-passthrough/vpn-passthrough@.service
[Unit]
Description=VPN Passthrough for %I

[Service]
Type=oneshot
ExecStart=true
RemainsAfterExit=yes