# vpn-passthrough/vpn-passthrough/vpn-passthrough-veth@.service
[Unit]
Description=Virtual Ethernet for VPN Passthrough for %I
BindTo=vpn-passthrough-netns@%i.service
After=vpn-passthrough-netns@%i.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=vpn-passthrough --app "%I" veth create
ExecStop=vpn-passthrough --app "%I" veth destroy